# debian/postinst
#!/bin/sh
set -e

case "$1" in
	configure)
		if ! getent group reverse-bin >/dev/null; then
			addgroup --system reverse-bin
		fi
		if ! getent passwd reverse-bin >/dev/null; then
			adduser --system --ingroup reverse-bin \
				--home /var/lib/reverse-bin/home --no-create-home \
				--shell /usr/sbin/nologin reverse-bin
		fi

		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin
		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin/apps
		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin/apps/logs
		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin/apps/logs/caddy-logs
		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin/apps/skills/reverse-bin-web-apps
		install -d -o reverse-bin -g reverse-bin /var/lib/reverse-bin/home
		install -d -m 0755 -o reverse-bin -g reverse-bin /var/lib/reverse-bin/keys

		if [ ! -f /var/lib/reverse-bin/keys/age.key ]; then
			/usr/lib/reverse-bin/age-keygen -o /var/lib/reverse-bin/keys/age.key
			/usr/lib/reverse-bin/age-keygen -y /var/lib/reverse-bin/keys/age.key > /var/lib/reverse-bin/keys/age.pub
		elif [ ! -f /var/lib/reverse-bin/keys/age.pub ]; then
			/usr/lib/reverse-bin/age-keygen -y /var/lib/reverse-bin/keys/age.key > /var/lib/reverse-bin/keys/age.pub
		fi
		chown reverse-bin:reverse-bin /var/lib/reverse-bin/keys/age.key
		chown reverse-bin:reverse-bin /var/lib/reverse-bin/keys/age.pub
		chmod 600 /var/lib/reverse-bin/keys/age.key
		chmod 644 /var/lib/reverse-bin/keys/age.pub

		if [ -d /usr/share/reverse-bin/logs-app ]; then
			for file in .env Caddyfile setup.sh README.md; do
				if [ ! -e "/var/lib/reverse-bin/apps/logs/$file" ]; then
					install -m 0644 -o reverse-bin -g reverse-bin "/usr/share/reverse-bin/logs-app/$file" "/var/lib/reverse-bin/apps/logs/$file"
				fi
			done
			chmod 0755 /var/lib/reverse-bin/apps/logs/setup.sh
			chown reverse-bin:reverse-bin /var/lib/reverse-bin/apps/logs/setup.sh
		fi

		if [ -f /usr/share/reverse-bin/skills/reverse-bin-web-apps/SKILL.md ]; then
			install -m 0644 -o reverse-bin -g reverse-bin /usr/share/reverse-bin/skills/reverse-bin-web-apps/SKILL.md /var/lib/reverse-bin/apps/skills/reverse-bin-web-apps/SKILL.md
		fi

		if [ ! -f /var/lib/reverse-bin/apps/.sops.yaml ]; then
			cat > /var/lib/reverse-bin/apps/.sops.yaml <<'EOF'
# See /usr/share/doc/reverse-bin/README.md.gz section "Encrypted app env files".
creation_rules:
  - path_regex: secrets\.enc\.json$
    key_groups:
      - age:
EOF
			printf '        # reverse-bin server key: /var/lib/reverse-bin/keys/age.pub\n' >> /var/lib/reverse-bin/apps/.sops.yaml
			sed 's/^/        - /' /var/lib/reverse-bin/keys/age.pub >> /var/lib/reverse-bin/apps/.sops.yaml
			chown reverse-bin:reverse-bin /var/lib/reverse-bin/apps/.sops.yaml
			chmod 644 /var/lib/reverse-bin/apps/.sops.yaml
		fi

		setcap 'cap_net_bind_service=+ep' /usr/bin/reverse-bin-caddy

		if command -v systemctl >/dev/null; then
			systemctl daemon-reload
			if systemctl is-enabled --quiet reverse-bin.service; then
				systemctl restart reverse-bin.service
			fi
		fi
		;;
esac

#DEBHELPER#
